# software/sorteros/sorteros-firstboot.service
[Unit]
Description=SorterOS first-boot setup (regen host keys, recreate swap, grow rootfs)
After=local-fs.target
Before=ssh.service sshd.service network-online.target
ConditionPathExists=!/var/lib/sorteros/firstboot-done

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/sorteros-firstboot.sh
RemainAfterExit=yes
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
